# .crontab.sample
*/5 * * * * cd congestion && make save_cloud
